# feh-go.cron
2 * * * *	root /etc/feh-go/feh-go
0 14 10 * *	root /etc/feh-go/feh-go upload
0 9 11 * *	root /etc/feh-go/feh-go backup
