# cmd/scripts/feh.cron
3 * * * *	root /etc/feh/feh
0 14 10 * *	root /etc/feh/feh upload
0 9 11 * *	root /etc/feh/feh backup
